# ansible/roles/notify/templates/inotify.service.j2
[Unit]
Description=inotifywait script for {{ item.item.log }}

[Service]
ExecStart=/tmp/{{ item.item.dest }}.sh

[Install]
WantedBy=multi-user.target